# scheduler/etc/cron.d/scheduler_commands
*/1 * * * *    root    flock -n /tmp/scheduler_run_subtasks.lock -c "scheduler run_subtasks"
