# scheduler/etc/supervisord.conf.d/tasktiger.conf
[program:tasktiger]
command=/usr/local/bin/tasktiger
